# .config/systemd/user/bazel-garbage-collect.service
[Unit]
Description=Remove dangling Bazel workspaces

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bazel-garbage-collect
IOSchedulingClass=idle
