# .config/systemd/user/bazel-garbage-collect.timer
[Unit]
Description=Remove dangling Bazel workspaces

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
