# .config/systemd/user/nix-collect-garbage.timer
[Unit]
Description=Delete unreachable nix store paths

[Timer]
OnStartupSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
