# .config/systemd/user/tmp-cleanup.timer
[Unit]
Description=Cleanup of Temporary Directories

[Timer]
OnStartupSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
