# dot_config/systemd/user/emacs.service.d/override.conf
[Service]
Environment=DOOMDIR=%h/.local/share/chezmoi/doom
Environment=DOOMLOCALDIR=%h/.local/share/doom
Environment=SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh
Restart=Always
