# dot_config/systemd/user/emacs.service.d/override.conf
[Service]
Environment=DOOMDIR=%h/.local/share/chezmoi/doom
ExecStart=
ExecStart=/usr/bin/emacs --fg-daemon
TimeoutSec=900
