# dot_config/systemd/user/emacs.service.d/override.conf
[Service]
Environment=DOOMDIR=~/.local/share/chezmoi/doom
TimeoutSec=900
