# dot_config/systemd/user/emacs.service.d/override.conf
[Service]
Environment=DOOMDIR=%h/.local/share/chezmoi/doom
Environment=SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh
ExecStart=
ExecStart=/usr/bin/emacs --fg-daemon
TimeoutSec=900
