# .config/systemd/user/emacs.service.d/override.conf
[Service]
Environment=SSH_AUTH_SOCK=%t/keyring/ssh
Restart=always
